# mt-cc
#a/!/bin/bash
#
# chkconfig: 12345 25 75
# description: Controls availability of block devices
#
# For Red-Hat-based distributions such as Fedora, RHEL, CentOS.
#	       
### BEGIN INIT INFO
# Provides: mt-cc
# Required-Start:
# Required-Stop:
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: Availability of block devices
### END INIT INFO

. /etc/init.d/functions

sbindir=/sbin

MTCC=mt-cc
MTCC_PATH="/opt/cs"
MTCC_BIN=$MTCC_PATH/$MTCC
MTCC_CONF=$MTCC_PATH/mt-cc.toml
LOGLEVEL=daemon.info

case "$1" in
  start)
	echo -n "Starting $MTCC service: "
	if [ ! -f $MTCC_CONF ]; then
		echo "Failed to find configuration file $MTCC_CONF"
		failure
		exit 1
	fi
	cd $MTCC_PATH
	nohup $MTCC_BIN $MTCC_CONF | logger -n localhost -p $LOGLEVEL -t $MTCC &
	RETVAL=$?
	if [ $RETVAL -eq 0 ]; then
		echo_success
		PID=$(pidofproc $MTCC)
		/usr/bin/logger -p $LOGLEVEL -n localhost -t "${MTCC}[${PID}]" started
	else
		failure
	fi
	;;

  stop)
	echo -n "Shutting down ${MTCC}: "
	PID=$(pidofproc $MTCC)
	killproc $MTCC
	[ "x$PID" != "x" ] && /usr/bin/logger -n localhost -p $LOGLEVEL -t "${MTCC}[${PID}]" stopped
	;;

  status)
	status $MTCC
	;;
  *)
	echo $"Usage: $0 {start|stop|status}"
	;;
esac
